# ssh.conf
[program:sshd]
command=sshd
priority=6920
stdout_logfile=/tmp/ssh_out
stdout_logfile_maxbytes=0
stderr_logfile=/dev/ssh_err
stderr_logfile_maxbytes=0
autorestart=true
